# deploy/systemd/obecnosc-backup.timer
[Unit]
Description=Codzienny backup aplikacji Obecnosc

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
